# update-nfs.service
[Unit]
Description=NFS DNS Updater

[Service]
ExecStart=/home/agoetz/bin/update-nfs.py